# scripts/drops-pipeline.cron
*/2 * * * * /var/www/html/gcz/scripts/drops-pipeline.sh
